# clean-home-dir.timer
[Unit]
Description=Clean Home Directory Daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
WantedBy=multi-user.target
